# ops/systemd/deribit_rr25.service
[Unit]
Description=Alpha12_24 Deribit RR25 Collector (public API)
After=network-online.target

[Service]
WorkingDirectory=%h/alpha12_24
Environment=PYTHONPATH=%h/alpha12_24
Environment=RUNS_DIR=%h/alpha12_24/runs
ExecStart=%h/alpha12_24/.venv/bin/python -m src.data.deribit_rr25 --currency BTC --interval_sec 60
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
